# scripts/systemd/decorai-watchdog.service
[Unit]
Description=DecorAI health monitor (restart app on /api/health failure)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# 按实际部署用户与路径修改以下三行：
User=admin
Group=admin
WorkingDirectory=/home/admin/home_sticker_design_api

Environment=DECORAI_HOME=/home/admin/home_sticker_design_api
Environment=DECORAI_WATCHDOG_LOG=/home/admin/home_sticker_design_api/logs/monitor/watchdog.log
Environment=APP_PORT=3000
Environment=WATCHDOG_INTERVAL=30
Environment=WATCHDOG_MIN_RESTART_GAP=120
Environment=WATCHDOG_CONSEC_FAILS=3
# 可选：经本机 Nginx 的入口自检（需与站点一致）
# Environment=DECORAI_EDGE_URL=http://127.0.0.1/
# 可选：rootless Podman 不稳时与 fix-server-502 一致
# Environment=DECORAI_USE_PODMAN_HOST_COMPOSE=1
# 可选：容器起来后重载 Nginx（需 sudoers 免密或 root 跑服务）
# Environment=DECORAI_RELOAD_NGINX_CMD=sudo /usr/bin/systemctl reload nginx

ExecStart=/home/admin/home_sticker_design_api/scripts/decorai-watchdog.sh
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target
